# web-scraper-docker/scripts/docker-host/cron.d/docker-jobs
3 0 * * * root /home/ec2-user/jobs/job-scraper-db.sh
3 1 * * * root /home/ec2-user/jobs/job-scraper-app.sh

